# packaging/quadlet/plottyd.container
[Unit]
Description=ploTTY daemon (container)

[Container]
Image=ghcr.io/yourorg/plotty:latest
Network=host
Volume=%h/.local/share/plotty:/var/lib/plotty:Z
Volume=%h/.config/plotty:/etc/plotty:Z
Environment=PLOTTY_CONFIG=/etc/plotty/config.yaml
Exec=plotty daemon --foreground

[Install]
WantedBy=default.target
